# directus-php-fpm/etc/supervisor.d/caddy.ini
[program:caddy]
directory=/opt/caddy/
command=/opt/caddy/caddy_linux_amd64 -conf /etc/caddy/Caddyfile
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/caddy_stdout.log
stderr_logfile=/var/log/supervisor/caddy_stderr.log
