# directus-php-fpm/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command = /usr/local/bin/run-directus
user = www-data
autostart = true
startsecs = 60
stdout_logfile=/var/log/supervisor/php-fpm_stdout.log
stderr_logfile=/var/log/supervisor/php-fpm_stderr.log
