# directus-php-fpm/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command = /usr/bin/php-fpm -F
user = www-data
autostart = true
stdout_logfile=/var/log/supervisor/php-fpm_stdout.log
stderr_logfile=/var/log/supervisor/php-fpm_stderr.log
